# archlive/airootfs/usr/local/share/noros/chroot/etc/systemd/system/update-issue.service
[Unit]
Description=Update /etc/issue with current IP addresses
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-issue

[Install]
WantedBy=multi-user.target
