# archlive/airootfs/usr/local/share/noros/chroot/etc/systemd/system/update-issue.timer
[Unit]
Description=Periodically update /etc/issue with current IP addresses

[Timer]
OnBootSec=5sec
OnUnitActiveSec=30sec

[Install]
WantedBy=timers.target
